# ai-service/systemd/ringrift-elo-sync.service
[Unit]
Description=RingRift Elo Sync Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ringrift/ai-service
Environment=PYTHONPATH=/home/ubuntu/ringrift/ai-service
ExecStart=/home/ubuntu/ringrift/ai-service/venv/bin/python scripts/elo_db_sync.py --mode worker --coordinator 100.77.77.122
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
