# ai-service/deploy/systemd/ringrift-keepalive.service
[Unit]
Description=RingRift Universal Keepalive Daemon
Documentation=https://github.com/an0mium/RingRift
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ringrift/ai-service
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Auto-detect node-id from hostname
ExecStart=/bin/bash -c 'exec python3 scripts/universal_keepalive.py --node-id "$(hostname)" --daemon'

# Restart policy - always restart on failure
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/ringrift/keepalive.log
StandardError=append:/var/log/ringrift/keepalive.log

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
