# ai-service/deploy/systemd/ringrift-resilience.service
[Unit]
Description=RingRift Node Resilience Daemon
After=network-online.target ringrift-p2p.service
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/ringrift/node.conf
WorkingDirectory=${RINGRIFT_DIR}
Environment=PYTHONPATH=${RINGRIFT_DIR}
ExecStart=/usr/bin/python3 ${RINGRIFT_DIR}/scripts/node_resilience.py --node-id ${NODE_ID} --coordinator ${COORDINATOR_URL} --ai-service-dir ${RINGRIFT_DIR} --p2p-port ${P2P_PORT}
Restart=always
RestartSec=30
StandardOutput=append:/var/log/ringrift/resilience.log
StandardError=append:/var/log/ringrift/resilience.log

[Install]
WantedBy=multi-user.target
